# javad.service.d/1.conf
[Service]
# don’t share /tmp with other processes
PrivateTmp=yes
