# javad.service.d/183.conf
[Service]
# no network access
PrivateNetwork=yes
