# javad.service.d/187.conf
[Service]
# never setuid or otherwise gain privileges
NoNewPrivileges=yes
# the system call filter: reset the filter to empty, then each subsequent assignment adds to it
SystemCallFilter=
# part of @basic-io
SystemCallFilter=read lseek
# part of @file-system (systemd commit 1a1b13c957, not in any release yet)
SystemCallFilter=open close stat fstat lstat getcwd access mmap readlink
# JIT code generation
SystemCallFilter=mprotect
# signals
SystemCallFilter=rt_sigaction
# threads
SystemCallFilter=clone futex sched_getaffinity arch_prctl
# required by nss-mymachines when DynamicUser=yes is in effect
SystemCallFilter=rt_sigprocmask
# don’t kill the process when an illegal syscall is issued, just return Operation not permitted
SystemCallErrorNumber=EPERM
