# javad.service.d/209.conf
[Service]
# no 32-bit binaries on 64-bit systems or other funny syscalls for foreign architectures
SystemCallArchitectures=native
# hide the real devices file system
PrivateDevices=yes
