# javad.service.d/211.conf
[Service]
# don’t talk to the internet (or other address families)
RestrictAddressFamilies=
