# javad.service.d/214.conf
[Service]
# make /usr, /boot and /etc read-only
# (superseded by ProtectSystem=strict on systemd ≥ v232)
ProtectSystem=full
# hide the home directories
ProtectHome=yes
