# javad.service.d/231.conf
[Service]
# don’t hog the CPU
RestrictRealtime=yes
