# javad.service.d/232.conf
[Service]
# don’t run as root
DynamicUser=yes
# hide other users
PrivateUsers=yes
# make the whole file system read-only
ProtectSystem=strict
# don’t write tunable kernel variables
ProtectKernelTunables=yes
# don’t edit cgroups
ProtectControlGroups=yes
# don’t explicitly load or unload kernel modules
ProtectKernelModules=yes
