# javad.service.d/no-user-namespaces.conf
[Service]
# use this on systems without user namespace support in the kernel (i.e. no CONFIG_USER_NS) (e.g. Arch Linux, see https://bugs.archlinux.org/task/36969)
PrivateUsers=no
