# javad.service.d/pre-214.conf
# ProtectSystem= and ProtectHome= were introduced in systemd v214;
# ProtectSystem=strict is emulated in pre-232.conf,
# emulate ProtectHome=yes for older versions

[Service]
ReadWriteDirectories=/home /run/user /root
