# javad.service.d/pre-232.conf
[Service]

# emulate ProtectSystem=strict + ProtectKernelTunables=yes + ProtectHome=yes
# see also systemd src/core/namespace.c
# (also, in systemd before v231, the directives are called *Directories instead of *Paths;
# systemd 231 still supports the old names, so we can safely use them in this file)

# make whole system read-only
ReadOnlyDirectories=/
# and completely hide home directories
InaccessibleDirectories=/home /run/user /root
# but exclude /dev, which is handled by PrivateDevices=
ReadWriteDirectories=/dev

# without PrivateUser=yes, ensure that we’re not running as root by specifying another user
User=nobody
