# javad.service.d/pre-233.conf
# Prior to systemd v233 (more specifically: before #5309 fixed #5308),
# ReadOnlyDirectories affects ExecStartPre= commands even with PermissionsStartOnly=yes,
# which means that writing and compiling the .java files fails.
# To fix this, explicitly make the working directory writable again.

[Service]
ReadWriteDirectories=/run/javad
