# javad.service.d/systemd-pre-232.conf
[Service]
# ProtectSystem=strict was introduced in systemd v232; emulate it for older versions
# based on systemd src/core/namespace.c, protect_system_strict_table
# make whole system read-only
ReadOnlyPaths=/
# except for API file systems, which are protected by
# ProtectKernelTunables=
ReadWritePaths=/proc /sys
# or PrivateDevices=
ReadWritePaths=/dev
# and home directories, which are protected by ProtectHome=
ReadWritePaths=/home /run/user /root
