# javad.service.d/systemd-pre-232.conf
# ProtectSystem=strict was introduced in systemd v232; emulate it for older versions
# based on systemd src/core/namespace.c, protect_system_strict_table
# (also, in systemd before v231, the directives are called *Directories instead of *Paths;
# systemd 231 still supports the old names, so we can safely use them in this file)

[Service]
# make whole system read-only
ReadOnlyDirectories=/
# except for API file systems, which are protected by
# ProtectKernelTunables=
ReadWriteDirectories=/proc /sys
# or PrivateDevices=
ReadWriteDirectories=/dev
# and home directories, which are protected by ProtectHome=
ReadWriteDirectories=/home /run/user /root
